# x86_64/packages/athena-welcome-app/etc/skel/.config/autostart/athena-welcome-app.desktop
[Desktop Entry]
Name=Athena Welcome App
GenericName=Athena Welcome App
X-GNOME-FullName=Athena Welcome App
Comment=Welcome Application for Athena Installations
Exec=/usr/local/bin/athena-welcome-app
Icon=arcolinux-hello
Terminal=false
Type=Application
Categories=GTK;Settings;Security;X-GNOME-Settings-Panel;X-GNOME-SystemSettings;X-Unity-Settings-Panel;X-XFCE-SettingsDialog;X-XFCE-SystemSettings;
Keywords=welcome;tool;
StartupNotify=true
